# luci-app-echworkers/root/etc/init.d/echworkers
#!/bin/sh /etc/rc.common

START=99
STOP=10
USE_PROCD=1

PROG=/usr/bin/ech-workers
CONFIG=echworkers
NFT_TABLE=echworkers
DNS_PORT=5335
DNSMASQ_CONF=/etc/dnsmasq.conf
GFW_MARKER="# ECH-Workers GFW DNS"
GFWLIST_FILE="/etc/echworkers/gfwlist.conf"

# 设置 nftables 透明代理规则（直接重定向到 TPROXY 端口）
setup_nftables() {
	local tproxy_port="$1"
	local server_ip="$2"
	
	# 清理旧规则
	nft delete table inet $NFT_TABLE 2>/dev/null
	
	# 创建新表和链
	nft add table inet $NFT_TABLE
	nft add chain inet $NFT_TABLE prerouting { type nat hook prerouting priority -100 \; }
	
	# === IPv4 规则 ===
	# 排除本地地址、私有地址
	nft add rule inet $NFT_TABLE prerouting ip daddr 127.0.0.0/8 return
	nft add rule inet $NFT_TABLE prerouting ip daddr 10.0.0.0/8 return
	nft add rule inet $NFT_TABLE prerouting ip daddr 172.16.0.0/12 return
	nft add rule inet $NFT_TABLE prerouting ip daddr 192.168.0.0/16 return
	nft add rule inet $NFT_TABLE prerouting ip daddr 224.0.0.0/4 return
	nft add rule inet $NFT_TABLE prerouting ip daddr 240.0.0.0/4 return
	
	# 排除代理服务器 IP
	if [ -n "$server_ip" ]; then
		local resolved_ip=$(echo "$server_ip" | grep -E '^[0-9]+\.[0-9]+\.[0-9]+\.[0-9]+$')
		if [ -z "$resolved_ip" ]; then
			resolved_ip=$(nslookup "$server_ip" 2>/dev/null | grep -A1 'Name:' | grep 'Address' | head -1 | awk '{print $2}')
		fi
		[ -n "$resolved_ip" ] && nft add rule inet $NFT_TABLE prerouting ip daddr "$resolved_ip" return
	fi
	
	# 重定向来自 LAN 的 IPv4 HTTP/HTTPS 流量到 TPROXY 端口
	nft add rule inet $NFT_TABLE prerouting iifname "br-lan" ip version 4 tcp dport 80 redirect to :$tproxy_port
	nft add rule inet $NFT_TABLE prerouting iifname "br-lan" ip version 4 tcp dport 443 redirect to :$tproxy_port
	
	# === IPv6 规则 ===
	# 排除 IPv6 本地地址
	nft add rule inet $NFT_TABLE prerouting ip6 daddr ::1/128 return
	nft add rule inet $NFT_TABLE prerouting ip6 daddr fe80::/10 return
	nft add rule inet $NFT_TABLE prerouting ip6 daddr fc00::/7 return
	nft add rule inet $NFT_TABLE prerouting ip6 daddr ff00::/8 return
	
	# 暂时阻止 IPv6 的 HTTP/HTTPS（因为 TPROXY 暂不支持 IPv6）
	# 这会强制客户端回退到 IPv4
	nft add rule inet $NFT_TABLE prerouting iifname "br-lan" ip6 version 6 tcp dport 80 drop
	nft add rule inet $NFT_TABLE prerouting iifname "br-lan" ip6 version 6 tcp dport 443 drop
	
	logger -t echworkers "nftables TPROXY rules applied (port: $tproxy_port, IPv6 blocked)"
}

# 清理 nftables 规则
cleanup_nftables() {
	nft delete table inet $NFT_TABLE 2>/dev/null
	logger -t echworkers "nftables rules removed"
}

# 设置 DNS 防污染（dnsmasq 分流 + dns2socks）
setup_dns() {
	local listen_port="$1"
	
	# 停止可能存在的进程
	killall dns2socks 2>/dev/null
	
	# 检查依赖
	if [ ! -x /usr/bin/dns2socks ]; then
		logger -t echworkers "dns2socks not found, DNS anti-pollution disabled"
		return 1
	fi
	
	# 启动 dns2socks（通过 SOCKS5 代理查询 8.8.8.8）
	dns2socks 127.0.0.1:$listen_port 8.8.8.8:53 127.0.0.1:$DNS_PORT /q &
	sleep 1
	
	# 配置 dnsmasq：GFW 域名走 dns2socks，其他走默认
	sed -i "/$GFW_MARKER/,\$d" $DNSMASQ_CONF
	
	# 添加 GFW 列表配置
	echo "" >> $DNSMASQ_CONF
	echo "$GFW_MARKER" >> $DNSMASQ_CONF
	echo "filter-AAAA" >> $DNSMASQ_CONF
	
	# 如果 GFW 列表文件存在，使用它；否则使用内置列表
	if [ -f "$GFWLIST_FILE" ] && [ -s "$GFWLIST_FILE" ]; then
		grep "^server=/" "$GFWLIST_FILE" >> $DNSMASQ_CONF
		local count=$(grep -c "^server=/" "$GFWLIST_FILE")
		logger -t echworkers "DNS anti-pollution enabled (dns2socks:$DNS_PORT, $count domains from gfwlist)"
	else
		# 内置基础列表
		cat >> $DNSMASQ_CONF << EOF
server=/google.com/127.0.0.1#$DNS_PORT
server=/youtube.com/127.0.0.1#$DNS_PORT
server=/googleapis.com/127.0.0.1#$DNS_PORT
server=/gstatic.com/127.0.0.1#$DNS_PORT
server=/googlevideo.com/127.0.0.1#$DNS_PORT
server=/ytimg.com/127.0.0.1#$DNS_PORT
server=/ggpht.com/127.0.0.1#$DNS_PORT
server=/twitter.com/127.0.0.1#$DNS_PORT
server=/x.com/127.0.0.1#$DNS_PORT
server=/facebook.com/127.0.0.1#$DNS_PORT
server=/instagram.com/127.0.0.1#$DNS_PORT
server=/telegram.org/127.0.0.1#$DNS_PORT
server=/t.me/127.0.0.1#$DNS_PORT
server=/github.com/127.0.0.1#$DNS_PORT
server=/githubusercontent.com/127.0.0.1#$DNS_PORT
server=/whatsapp.com/127.0.0.1#$DNS_PORT
server=/signal.org/127.0.0.1#$DNS_PORT
server=/openai.com/127.0.0.1#$DNS_PORT
server=/anthropic.com/127.0.0.1#$DNS_PORT
server=/claude.ai/127.0.0.1#$DNS_PORT
server=/chatgpt.com/127.0.0.1#$DNS_PORT
server=/wikipedia.org/127.0.0.1#$DNS_PORT
server=/wikimedia.org/127.0.0.1#$DNS_PORT
server=/reddit.com/127.0.0.1#$DNS_PORT
server=/twitch.tv/127.0.0.1#$DNS_PORT
server=/discord.com/127.0.0.1#$DNS_PORT
server=/spotify.com/127.0.0.1#$DNS_PORT
server=/netflix.com/127.0.0.1#$DNS_PORT
EOF
		logger -t echworkers "DNS anti-pollution enabled (dns2socks:$DNS_PORT, builtin list)"
	fi
	
	# 异步重启 dnsmasq，避免阻塞
	( sleep 1; /etc/init.d/dnsmasq restart ) &
}

# 清理 DNS 配置
cleanup_dns() {
	killall dns2socks 2>/dev/null
	
	# 移除 dnsmasq 配置
	sed -i "/$GFW_MARKER/,\$d" $DNSMASQ_CONF
	
	# 异步重启 dnsmasq，避免阻塞
	( sleep 1; /etc/init.d/dnsmasq restart ) &
	
	logger -t echworkers "DNS anti-pollution disabled"
}

start_service() {
	config_load "$CONFIG"
	
	local enabled server listen token ip dns ech routing transparent tproxy_port
	config_get enabled config enabled 0
	config_get server config server
	config_get listen config listen "127.0.0.1:30001"
	config_get token config token
	config_get ip config ip
	config_get dns config dns
	config_get ech config ech
	config_get routing config routing "bypass_cn"
	config_get transparent config transparent 0
	config_get tproxy_port config tproxy_port 12345
	
	[ "$enabled" = "1" ] || return 0
	[ -n "$server" ] || return 1
	[ -x "$PROG" ] || {
		logger -t echworkers "ech-workers binary not found at $PROG"
		return 1
	}
	
	# 构建命令参数
	local args="-f $server -l $listen"
	[ -n "$token" ] && args="$args -token $token"
	[ -n "$ip" ] && args="$args -ip $ip"
	[ -n "$dns" ] && args="$args -dns $dns"
	[ -n "$ech" ] && args="$args -ech $ech"
	[ -n "$routing" ] && args="$args -routing $routing"
	
	# 透明代理模式：添加 TPROXY 参数
	if [ "$transparent" = "1" ]; then
		args="$args -tproxy 0.0.0.0:$tproxy_port"
	fi
	
	# 启动 ech-workers
	procd_open_instance echworkers
	procd_set_param command /bin/sh -c "cd /usr/share/ech-workers && exec $PROG $args"
	procd_set_param respawn
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_set_param file /etc/config/echworkers
	procd_close_instance
	
	# 透明代理模式：异步设置 nftables 规则和 DNS 防污染
	if [ "$transparent" = "1" ]; then
		(
			# 等待服务启动
			sleep 2
			
			# 解析服务器 IP（用于排除规则）
			local server_ip_resolved=""
			if [ -n "$ip" ]; then
				server_ip_resolved="$ip"
			else
				local server_host=$(echo "$server" | cut -d: -f1)
				server_ip_resolved=$(nslookup "$server_host" 2>/dev/null | grep -A1 'Name:' | grep 'Address' | head -1 | awk '{print $2}')
			fi
			
			setup_nftables "$tproxy_port" "$server_ip_resolved"
			
			# 设置 DNS 防污染
			local listen_port=$(echo "$listen" | cut -d: -f2)
			setup_dns "$listen_port"
			
			# 设置定时更新
			setup_cron
		) &
	fi
}

stop_service() {
	# 强制停止所有相关进程
	killall -9 ech-workers 2>/dev/null
	killall -9 dns2socks 2>/dev/null
	
	# 清理 nftables 和 cron
	cleanup_nftables
	cleanup_cron
	
	# 清理 dnsmasq 配置（不重启，避免网络中断）
	sed -i "/$GFW_MARKER/,\$d" $DNSMASQ_CONF
	
	# 使用 dnsmasq reload 而不是 restart（不会中断网络）
	killall -HUP dnsmasq 2>/dev/null
	
	logger -t echworkers "Service stopped"
}

service_triggers() {
	procd_add_reload_trigger "$CONFIG"
}

reload_service() {
	# 使用 nohup 完全脱离，避免阻塞 procd
	nohup /usr/bin/echworkers-control restart >/dev/null 2>&1 &
}

# 设置透明代理（供 echworkers-control 调用）
setup_transparent() {
	config_load "$CONFIG"
	
	local server listen ip transparent tproxy_port
	config_get server config server
	config_get listen config listen "127.0.0.1:30001"
	config_get ip config ip
	config_get transparent config transparent 0
	config_get tproxy_port config tproxy_port 12345
	
	[ "$transparent" = "1" ] || return 0
	
	# 解析服务器 IP
	local server_ip_resolved=""
	if [ -n "$ip" ]; then
		server_ip_resolved="$ip"
	else
		local server_host=$(echo "$server" | cut -d: -f1)
		server_ip_resolved=$(nslookup "$server_host" 2>/dev/null | grep -A1 'Name:' | grep 'Address' | head -1 | awk '{print $2}')
	fi
	
	setup_nftables "$tproxy_port" "$server_ip_resolved"
	
	local listen_port=$(echo "$listen" | cut -d: -f2)
	setup_dns "$listen_port"
	
	setup_cron
}

# 设置定时更新 GFW 列表
setup_cron() {
	local auto_update interval
	config_get auto_update gfwlist auto_update 0
	config_get interval gfwlist update_interval 24
	
	# 移除旧的 cron 任务
	sed -i '/echworkers-gfwlist/d' /etc/crontabs/root 2>/dev/null
	
	if [ "$auto_update" = "1" ]; then
		# 根据间隔设置 cron
		case "$interval" in
			6)  echo "0 */6 * * * /usr/bin/echworkers-gfwlist update >/dev/null 2>&1" >> /etc/crontabs/root ;;
			12) echo "0 */12 * * * /usr/bin/echworkers-gfwlist update >/dev/null 2>&1" >> /etc/crontabs/root ;;
			24) echo "0 4 * * * /usr/bin/echworkers-gfwlist update >/dev/null 2>&1" >> /etc/crontabs/root ;;
			72) echo "0 4 */3 * * /usr/bin/echworkers-gfwlist update >/dev/null 2>&1" >> /etc/crontabs/root ;;
			168) echo "0 4 * * 0 /usr/bin/echworkers-gfwlist update >/dev/null 2>&1" >> /etc/crontabs/root ;;
		esac
		/etc/init.d/cron restart
		logger -t echworkers "GFW list auto-update enabled (every $interval hours)"
	fi
}

# 清理 cron 任务
cleanup_cron() {
	sed -i '/echworkers-gfwlist/d' /etc/crontabs/root 2>/dev/null
	/etc/init.d/cron restart 2>/dev/null
}
